# mush.service
[Unit]
Description=MushMeet.com
After=network-online.target networking.service
Wants=network-online.target networking.service

[Service]
Type=simple
WorkingDirectory=/home/jdoe/www/
ExecStart=uvicorn main:app --port 80 --host 0.0.0.0

[Install]
WantedBy=multi-user.target
